# bash/templates/supervisor/client.conf
; FreeCP — Supervisor: {{DOMAIN}}

[program:{{CONTAINER_NAME}}_queue]
command=docker exec {{CONTAINER_NAME}} php artisan queue:work redis
    --sleep=3 --tries=3 --max-time=3600
    --queue=default,high,low
process_name=%(program_name)s_%(process_num)02d
numprocs={{QUEUE_WORKERS}}
autostart=true
autorestart=true
stopwaitsecs=60
stderr_logfile=/opt/freecp/clients/{{DOMAIN}}/logs/queue.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/dev/null
priority=20

[program:{{CONTAINER_NAME}}_scheduler]
command=/bin/bash -c "while true; do docker exec {{CONTAINER_NAME}} php artisan schedule:run --no-interaction; sleep 60; done"
autostart=true
autorestart=true
stderr_logfile=/opt/freecp/clients/{{DOMAIN}}/logs/scheduler.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/dev/null
priority=30

[group:{{CONTAINER_NAME}}]
programs={{CONTAINER_NAME}}_queue,{{CONTAINER_NAME}}_scheduler